# config/systemd/system/logid.service.d/release-keys.conf
[Service]
ExecStopPost=-/usr/local/bin/logid-release-keys.py
